# .docker/django/rootfs/etc/periodic/hourly/django
#!/usr/bin/with-contenv sh

cd /app && python /app/manage.py runjobs hourly
