# .docker/django/rootfs/var/spool/cron/crontabs/root
#!/usr/bin/with-contenv sh

03 * * * * cd /app && python /app/manage.py runjobs hourly >/proc/1/fd/1 2>/proc/1/fd/2
