# .docker/django/rootfs/etc/services.d/cron/run
#!/usr/bin/with-contenv sh

exec /usr/sbin/crond -f -l 8
